# elal-monitor.service
[Unit]
Description=El Al Flight Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=eylonc11092001
WorkingDirectory=/home/eylonc11092001/elal_bot
EnvironmentFile=/home/eylonc11092001/elal_bot/.env
ExecStart=/home/eylonc11092001/elal_bot/venv/bin/python monitor.py
Restart=on-failure
RestartSec=300
StartLimitIntervalSec=600
StartLimitBurst=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
